# etc/ntplogtemp.timer
[Unit]
Description=Run ntplogtemp every five minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=ntpd.service
