# etc/ntpviz-weekly.timer
[Unit]
Description=Run ntpviz for daily weekly graphing twice a day

[Timer]
OnCalendar=11,23:45:00
Persistent=true

[Install]
WantedBy=ntpd.service
